# image-background-remover.service
[Unit]
Description=AI Image Background Remover
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/image-background-remover
ExecStart=/usr/bin/gunicorn --bind 0.0.0.0:5000 app:app
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
